# systemd-templates/hardware-camera.service.template
[Unit]
Description=Fish Tank Camera Stream
After=network-online.target

[Service]
User={USERNAME}
# this assumes MediaMTX is running on the backend server.
ExecStart=/bin/bash -c 'rm -f /tmp/{CAM_NUMBER}_pipe.h264 && mkfifo /tmp/{CAM_NUMBER}_pipe.h264 && (/usr/bin/ffmpeg -use_wallclock_as_timestamps 1 -f h264 -i /tmp/{CAM_NUMBER}_pipe.h264 -c:v copy -rtsp_transport tcp -f rtsp -pkt_size 1200 {RTSP_URL}/{CAM_NUMBER} &) && /usr/bin/rpicam-vid -t 0 --inline --intra 30 --mode 2304:1296 --width 1280 --height 720 --profile high --level 4.2 --bitrate 5000000 --denoise cdn_off -o /tmp/{CAM_NUMBER}_pipe.h264'
ExecStopPost=/usr/bin/killall ffmpeg
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
